# root/etc/services.d/betanin/run
#!/usr/bin/with-contenv bash
umask 022

exec \
	s6-setuidgid abc betanin
